# scripts/smart-git-proxy.service
[Unit]
Description=Smart Git HTTP mirror proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=smart-git-proxy
Group=smart-git-proxy
EnvironmentFile=/etc/smart-git-proxy/env
ExecStart=/usr/bin/smart-git-proxy
Restart=on-failure
RestartSec=2
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
